# templates/default/cadvisor.service.erb
[Unit]
Description = cAdvisor
StartLimitIntervalSec=30
StartLimitBurst=20

[Service]
User = <%= node['ndb']['user'] %>
ExecStart = <%= @cadvisor_bin %> \
--port=:<%= node['hops']['cadvisor']['port'] %>
LimitNOFILE = 20000
Restart=on-failure

Type = simple

[Install]
WantedBy = multi-user.target